# supervisord.conf
# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log

[program:spring-boot-app]
command=java -jar /app/app.jar
autorestart=true
stderr_logfile=/var/log/spring-boot-app.err.log
stdout_logfile=/var/log/spring-boot-app.out.log
